# pkgbuilds/unity-base/unity/pkg/unity/usr/share/upstart/sessions/unity-panel-service-lockscreen.conf
description "Backing Service for the Unity Panel"
author "Andrea Azzarone <andrea.azzarone@canonical.com>"

start on desktop-lock
stop on desktop-unlock

pre-start exec /usr/lib/unity/upstart-prestart-check

respawn
exec /usr/lib/unity/unity-panel-service --lockscreen-mode
