# pkgbuilds/unity-indicators/indicator-sound/src/build/data/indicator-sound.service
[Unit]
Description=Indicator Sound Service
PartOf=graphical-session.target
After=indicators-pre.target

[Service]
ExecStart=/usr/libexec/indicator-sound/indicator-sound-service
Restart=on-failure
