# pkgbuilds/unity-base/unity/src/build/services/unity-panel-service.service
[Unit]
Description=Backing Service for the Unity Panel
After=unity7.service
PartOf=graphical-session.target
BindsTo=indicators-pre.target

[Service]
ExecStartPre=/usr/lib/unity/systemd-prestart-check
ExecStart=/usr/lib/unity/unity-panel-service
Restart=on-failure
