# pkgbuilds/unity-indicators/indicator-session/src/build/data/indicator-session.service
[Unit]
Description=Indicator Session Service
PartOf=graphical-session.target
After=indicators-pre.target

[Service]
ExecStart=/usr/libexec/indicator-session/indicator-session-service
Restart=on-failure
