# pkgbuilds/unity-base/hud/src/hud-14.10+17.10.20170619/data/hud.service.in
[Unit]
Description=Unity HUD
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=com.canonical.hud
ExecStart=@pkglibexecdir@/hud-service
