# pkgbuilds/unity-base/hud/src/hud-14.10+17.10.20170619/data/window-stack-bridge.service.in
[Unit]
Description=Window Stack Bridge
PartOf=graphical-session.target

[Service]
ExecStart=@pkglibexecdir@/window-stack-bridge
Restart=on-failure
