# pkgbuilds/unity-base/hud/pkg/hud/usr/lib/systemd/user/window-stack-bridge.service
[Unit]
Description=Window Stack Bridge
PartOf=graphical-session.target

[Service]
ExecStart=/usr/libexec/hud/window-stack-bridge
Restart=on-failure
